# docker/php/crontab
SHELL=/bin/bash
* * * * * cd /var/www && /usr/local/bin/php /var/www/artisan schedule:run >> /var/www/storage/logs/cron.log 2>&1
