# docker/php/supervisor.conf
[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock
file=/etc/supervisor/conf.d/supervisord.conf

[program:cron]
command = cron -f
autostart=true
autorestart=true
user=root

[program:php]
command = docker-php-entrypoint php-fpm
